# docker/s6-rc.d/init-db/up
#!/command/with-contenv sh
set -e
echo "[init-db] Running database migrations..."
node /app/migrate.js
echo "[init-db] Done."
